# supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir=/var/log/supervisor            ; ('AUTO' child log dir, default $TEMP)

nodaemon=true

[program:mariadb]
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe
autorestart=true

[program:firstrun]
command=/opt/firstrun.sh
startsecs = 0
autorestart=false

[program:shootmania]
command=/opt/speedball.sh
autorestart=true

[program:maniacontrol]
command=/usr/bin/pidproxy /opt/ManiaControl/ManiaControl.pid /opt/maniacontrol.sh
autorestart=true

; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or
; newlines).  It can also contain wildcards.  The filenames are
; interpreted as relative to this file.  Included files *cannot*
; include files themselves.

[include]
files = /etc/supervisor/conf.d/*.conf
